# assets/adlists.cron
# -*- mode: crontab; -*-
LOGFILE=/root/update-adlists.log

# MIN HOUR DAY MONTH DOW USER CMD
* * * * * root /opt/adlists/update-adlists.sh /etc/adlists.list /etc/pihole/gravity.db &>>"$LOGFILE"
